# Podman/quadlet/birdnet-go-autotls.container
[Unit]
Description=BirdNET-Go Container with AutoTLS
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/tphakala/birdnet-go:nightly
ContainerName=birdnet-go-autotls
Volume=./config:/config
Volume=./data:/data
Volume=./certs:/certs
# AutoTLS binds the ACME HTTP-01 + HTTP->HTTPS redirect listener on 8080 and
# HTTPS on 8443 inside the container, so map the privileged host ports to those
# non-privileged container ports (no NET_BIND_SERVICE needed). Mirrors
# Podman/podman-compose.autotls.yml.
PublishPort=80:8080
PublishPort=443:8443
Environment=TZ=UTC
Environment=BIRDNET_UID=1000
Environment=BIRDNET_GID=1000
Environment=BIRDNET_SECURITY_AUTOTLS=true
Environment=BIRDNET_SECURITY_HOST=%i
Device=/dev/snd:/dev/snd
Network=bridge
Tmpfs=/config/hls:exec,size=50M,uid=1000,gid=1000,mode=0755

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target